# profile.d/bigredbutton.sh
bigredbutton_activate() {
    source /var/www/html/bigredbutton-py/src/brb_env/bin/activate #"$@"
}
export -f bigredbutton_activate
